# pi-remote-daemon/deploy/launchd/dev.pi-remote.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- SPDX-License-Identifier: MIT -->
<!--
    launchd user agent template for pi-remote-daemon (SPEC.md § D11).
    Install at ~/Library/LaunchAgents/dev.pi-remote.daemon.plist and
    `launchctl load` it from the install script.

    Replace {{HOME}} with the absolute home directory of the daemon's owner
    when this template is materialized.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.pi-remote.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/pi-remote-daemon</string>
        <string>--config</string>
        <string>{{HOME}}/.config/pi-remote/daemon.toml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/.pi-remote/daemon.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/.pi-remote/daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
